# rootfs/etc/supervisor/conf.d/pcscd.conf
[program:pcscd]
command=/usr/sbin/pcscd --foreground --auto-exit
autostart=true
autorestart=true
stderr_logfile=/var/log/pcscd.err.log
stdout_logfile=/var/log/pcscd.out.log

